# pxe/common_files/check-scout-updates.timer
[Unit]
Description=Run scout daily update check

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
