# pxe/common_files/check-scout-updates.service
[Unit]
Description=Check for scout updates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/forge/check-scout-updates.sh
