# crates/agent/templates/default-isc-dhcp-relay_etv.conf
# Auto-generated by Forge!
[program: isc-dhcp-relay-default]
command = /usr/sbin/dhcrelay --nl -d {{range .Uplinks}} -i {{.}} {{end}} {{range .VlanIDs}} -i vlan{{.}} {{ end }} -U lo{{range .Servers}} {{.}}{{end}} -r {{.RemoteId}}
autostart = true
autorestart = unexpected
startsecs = 3
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 3
